# utils/91-havit-h2002yg.rules
ACTION=="add", SUBSYSTEM=="sound", KERNEL=="card*", ATTRS{idVendor}=="3206", ATTRS{idProduct}=="0881", RUN+="/usr/bin/amixer -c %n sset PCM 100%%"
